# scripts/deploy/omnibot-activity.service
[Unit]
Description=OmniBot Activity API and Frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omnibot
Group=omnibot
WorkingDirectory=/opt/omnibot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
EnvironmentFile=/opt/omnibot/.env
ExecStart=/opt/omnibot/.venv/bin/python -m uvicorn activity.server.main:app --host 127.0.0.1 --port 8008
Restart=always
RestartSec=5
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
ReadOnlyPaths=/opt/omnibot
ReadWritePaths=/opt/omnibot/data /opt/omnibot/logs

[Install]
WantedBy=multi-user.target
